# mcpproxy/com.mcpproxy.tray.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mcpproxy.tray</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/mcpproxy-tray</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/tim/Library/Logs/mcpproxy/tray-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/tim/Library/Logs/mcpproxy/tray-stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/tim</string>
        <key>SHELL</key>
        <string>/opt/homebrew/bin/fish</string>
    </dict>

    <key>ProcessType</key>
    <string>Interactive</string>

    <key>ExitTimeOut</key>
    <integer>30</integer>
</dict>
</plist>
